# platform/solaris/smf/opengrok.xml
<?xml version="1.0"?>
<!--
CDDL HEADER START

The contents of this file are subject to the terms of the
Common Development and Distribution License (the "License").
You may not use this file except in compliance with the License.

You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
or http://www.opensolaris.org/os/licensing.
See the License for the specific language governing permissions
and limitations under the License.

When distributing Covered Code, include this CDDL HEADER in each
file and include the License file at usr/src/OPENSOLARIS.LICENSE.
If applicable, add the following below this CDDL HEADER, with the
fields enclosed by brackets "[]" replaced with your own identifying
information: Portions Copyright [yyyy] [name of copyright owner]

CDDL HEADER END

Copyright 2010 Sun Microsystems, Inc. All rights reserved.
Use is subject to license terms.

-->
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="opengrok">
<service
	name="application/opengrok"
	type="service"
	version="1">

	<single_instance />

	<dependency
		name="filesystem"
		grouping="require_all"
		restart_on="none"
		type="service">
		<service_fmri value="svc:/system/filesystem/local:default" />
	</dependency>

	<dependency
		name="network"
		grouping="require_all"
		restart_on="none"
		type="service">
		<service_fmri value="svc:/milestone/network:default" />
	</dependency>

	<dependency
		name="autofs"
		grouping="require_all"
		restart_on="none"
		type="service">
		<service_fmri value="svc:/system/filesystem/autofs:default" />
	</dependency>

	<dependency
		name="tomcat"
		grouping="require_all"
		restart_on="none"
		type="service">
		<service_fmri value="svc:/network/http:tomcat6" />
	</dependency>

	<property_group
		name="opengrok"
		type="application">
		<propval
			name="maxmemory"
			type="astring"
			value="256" />
		<propval
			name="verbose"
			type="boolean"
			value="true" />
		<propval
			name="srcdir"
			type="astring"
			value="/var/opengrok/src" />
		<propval
			name="cachedir"
			type="astring"
			value="/var/opengrok/data" />
		<propval
			name="default_params"
			type="astring"
			value="-P -S -r on" />
		<propval
			name="java_extra_params"
			type="astring"
			value="" />
		<propval
			name="extra_params"
			type="astring"
			value="" />
		<propval
			name="readonly_config"
			type="astring"
			value="/etc/opengrok/readonly_configuration.xml" />
	</property_group>

	<instance
		name="default"
		enabled="false">
		<exec_method
			type="method"
			name="start"
			exec="/lib/svc/method/svc-opengrok start"
			timeout_seconds="30" >
			<method_context>
				<method_credential
					user="webservd"
					group="webservd" />
			</method_context>
		</exec_method>
		<exec_method
			type="method"
			name="stop"
			exec="/lib/svc/method/svc-opengrok stop %{restarter/contract}"
			timeout_seconds="60" >
			<method_context>
				<method_credential
					user="webservd"
					group="webservd" />
			</method_context>
		</exec_method>

		<exec_method
			type='method'
			name='refresh'
			exec='/lib/svc/method/svc-opengrok refresh'
			timeout_seconds='60' >
			<method_context>
				<method_credential
					user="webservd"
					group="webservd" />
			</method_context>
		</exec_method>
	</instance>

	<stability value="Unstable" />
	<template>
		<common_name>
			<loctext xml:lang="C">
			Opengrok application server
			</loctext>
		</common_name>
		<documentation>
			<manpage
				title="opengrok"
				section="1"
				manpath="/usr/opengrok/man" />
		</documentation>
	</template>
</service>
</service_bundle> 
