# platform/solaris/smf/opengrok.xml
<?xml version="1.0"?>
<!--
CDDL HEADER START

The contents of this file are subject to the terms of the
Common Development and Distribution License (the "License").
You may not use this file except in compliance with the License.

You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
or http://www.opensolaris.org/os/licensing.
See the License for the specific language governing permissions
and limitations under the License.

When distributing Covered Code, include this CDDL HEADER in each
file and include the License file at usr/src/OPENSOLARIS.LICENSE.
If applicable, add the following below this CDDL HEADER, with the
fields enclosed by brackets "[]" replaced with your own identifying
information: Portions Copyright [yyyy] [name of copyright owner]

CDDL HEADER END

Copyright 2008 Sun Microsystems, Inc. All rights reserved.
Use is subject to license terms.

-->
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="opengrok">
    <service name="application/servers/opengrok" type="service" version="1">
        <single_instance/>
        
        <dependency name="filesystem" grouping="require_all" restart_on="none" type="service">
            <service_fmri value="svc:/system/filesystem/local:default" />
        </dependency>

        <dependency name="network" grouping="require_all" restart_on="none" type="service">
            <service_fmri value="svc:/milestone/network:default" />
        </dependency>

        <property_group name="general" type="framework">
            <propval name="action_authorization" type="astring" value="solaris.smf.manage.opengrok" />
            <propval name="value_authorization" type="astring" value="solaris.smf.value.opengrok" />
        </property_group>
      
        <property_group name="opengrok" type="application">
            <propval name="java" type="astring" value="/usr/jdk/latest/bin/java"/>
            <propval name="java_options" type="astring" value="-server"/>
            <propval name="runtime_path" type="astring" value="/usr/bin:/usr/opengrok/bin"/>
            <propval name="configuration" type="astring" value="/etc/opengrok/opengrok.properties"/>
            <propval name="home" type="astring" value="/var/opengrok" />
            <propval name="value_authorization" type="astring"
                   value="solaris.smf.value.opengrok" />
            <propval name="action_authorization" type="astring"
                   value="solaris.smf.manage.opengrok" />
        </property_group>

        <instance name="opengrok" enabled="false">
            <exec_method type="method" name="start" exec="/lib/svc/method/opengrok start" timeout_seconds="30" >
                <method_context>
                    <method_credential user="noaccess" group="noaccess" />
                </method_context>
            </exec_method>
            <exec_method type="method" name="stop" exec="/lib/svc/method/opengrok start" timeout_seconds="60" >
                <method_context>
                    <method_credential user="noaccess" group="noaccess" />
                </method_context>
            </exec_method>
        </instance>

        <stability value="Unstable" />
        <template>
            <common_name>
                <loctext xml:lang="C"> Opengrok application server</loctext>
            </common_name>
            <documentation>
                <manpage title="opengrok" section="1" manpath="/usr/opengrok/man" />
            </documentation>
        </template>
    </service>
</service_bundle> 
